# systemd/sloef.service
[Unit]
Description=Check if web services are online

[Service]
Type=simple
ExecStart=/usr/bin/sloef

#[Install]
#WantedBy=multi-user.target
