# systemd/sloef.timer
[Unit]
Description=Check if web services are online

[Timer]
OnBootSec=1min
OnUnitActiveSec=2h

[Install]
WantedBy=timers.target
